# contact.socket
[Unit]
Description=gunicorn  contact socket

[Socket]
ListenStream=/run/ contact.socket

[Install]
WantedBy=sockets.target
